# crontab
0 7 * * * python /app/nhl_data_parser --update_db